# conf/udev/15-wifi-stick.rules
SUBSYSTEM=="usb",
ATTRS{idProduct}=="2573",
ATTRS{idVendor}=="148f",
RUN+="/usr/local/ipcam/bin/conn-health"
